# packages/ApplianceKit-NG/ApplianceKit-NG-0.6.1.tar.gz/ApplianceKit-NG-0.6.1/appliancekit/specs/arch/network.service.tmpl
[Unit]
Description=Static IP Connectivity
Wants=network.target
Before=network.target
BindsTo=sys-subsystem-net-devices-eth0.device
After=sys-subsystem-net-devices-eth0.device

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/ip link set dev eth0 up

{% if eth0 == 'dhcp' %}
ExecStart=/sbin/dhcpcd eth0
{% elif isinstance(eth0, dict) %}
ExecStart=/sbin/ip addr add {{ eth0['address'] }}/{{ eth0['netmask'] }} brd {{ eth0['broadcast'] }} dev eth0
ExecStart=/sbin/ip route add default via {{ eth0['gateway'] }}
{% endif %}

{% if eth0 == 'dhcp' %}
ExecStop=/sbin/dhcpcd -k eth0
{% endif %}

ExecStop=/sbin/ip addr flush dev eth0
ExecStop=/sbin/ip link set dev eth0 down

[Install]
WantedBy=multi-user.target
